# etc_scripts/naverpost-bot.service
# systemd service file for naverPost Telegram Bot (Enhanced DNS/Network Stability)
# Installation: sudo cp naverpost-bot.service /etc/systemd/system/
# Enable: sudo systemctl enable naverpost-bot
# Start: sudo systemctl start naverpost-bot

[Unit]
Description=NaverPost Telegram Bot Service with DNS Stability
# Wait for network to be fully ready (not just interface up)
After=network-online.target systemd-resolved.service
Wants=network-online.target
# Require network to be online (with working DNS)
Requires=network-online.target

[Service]
Type=simple
User=mini
Group=mini
WorkingDirectory=/home/mini/dev/naverPost
ExecStart=/home/mini/dev/naverPost/venv/bin/python3 /home/mini/dev/naverPost/etc_scripts/start_bot_with_health_check.py

# Environment
Environment=PYTHONPATH=/home/mini/dev/naverPost
# DNS fallback settings (override in .env if needed)
Environment=TELEGRAM_DNS_SERVERS=1.1.1.1,8.8.8.8,9.9.9.9
Environment=TELEGRAM_API_IPS=149.154.167.220,149.154.167.51,149.154.167.50
# Network readiness settings
Environment=NETWORK_WAIT_TIMEOUT=30
Environment=DNS_CHECK_ENABLED=true
# Playwright headless 모드 (서버 환경: 반드시 true)
Environment=HEADLESS=true
Environment=PLAYWRIGHT_HEADLESS=true
# Load additional environment from .env file
EnvironmentFile=-/home/mini/dev/naverPost/.env

# Enhanced process management for DNS stability
Restart=always
# Progressive backoff: start with 10s, then increase
RestartSec=10
# Allow more restart attempts but over longer period
StartLimitIntervalSec=600
StartLimitBurst=10
# Restart on any exit code (including DNS failures)
RestartForceExitStatus=1 2

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=naverpost-bot

# Security
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=false
ReadWritePaths=/home/mini/dev/naverPost/data
ReadWritePaths=/home/mini/dev/naverPost/uploads
ReadWritePaths=/home/mini/dev/naverPost/memory
ReadWritePaths=/home/mini/dev/naverPost/logs
PrivateTmp=yes

# Resource limits
MemoryMax=512M
TasksMax=100

# Graceful shutdown
TimeoutStopSec=30
KillMode=mixed
KillSignal=SIGTERM

[Install]
WantedBy=multi-user.target
